# osha@.timer
[Unit]
Description=Check factorio server updates for %i

[Timer]
OnUnitActiveSec=30m
Persistent=true

[Install]
WantedBy=timers.target
