# osha@.service
[Unit]
Description=Run factorio server updates for %i

[Service]
User=factorio
Type=oneshot
ExecStart=/usr/bin/osha %i
